# recipes-alaric/alaric-fpga-init/files/alaric-fpga-init.service
[Unit]
Description=FPGA Initialization

[Service]
Type=simple
ExecStart=/usr/bin/alaric-fpga-init.sh

[Install]
WantedBy=basic.target
